# systemd/yc-preemptible-healer.timer
[Unit]
Description=YC Preemptible Healer Timer

[Timer]
OnBootSec=60
OnUnitActiveSec=120

[Install]
WantedBy=timers.target
