# systemd/yc-preemptible-healer.service
[Unit]
Description=YC Preemptible Healer

[Service]
Type=oneshot
ExecStart=/opt/yc-preemptible-healer/heal.sh

