# systemd/yc-preemptible-healer.service
[Unit]
Description=YC Preemptible Healer
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/yc-preemptible-healer/heal.sh
